# deploy/systemd/open-ontologies.service
[Unit]
Description=OntoStar open-ontologies MCP server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/open-ontologies server serve_http --host 0.0.0.0 --port 3050
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
